# utils/jobs.cron

# Backup all data on live website
12 3 * * *		/var/www/www.inkscape.org/utils/backup

